# setup_panda/panda.cron
3 0 * * * ubuntu sudo -u panda -E python /opt/panda/manage.py purge_orphaned_uploads
